# pi-media-hub.service
[Unit]
Description=Pi Media Hub Kiosk
After=graphical.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%USERNAME%
Group=%USERNAME%
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/%USERNAME%/.Xauthority
WorkingDirectory=%INSTALL_DIR%
ExecStartPre=/bin/sleep 5
ExecStart=/usr/bin/python3 %INSTALL_DIR%/launcher.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=graphical.target
